# leonard/bin/leonard.service
[Unit]
Description="1951 Leonard Refrigerator Daemon"
After=network.target

[Service]
Environment="HOME=/home/root"
#User=root
ExecStart=/usr/bin/leonard
Restart=always
RestartSec=10s

