# leonard/etc/systemd/system/multi-user.target.wants/leonard-api.service
[Unit]
Description="Leonard REST API"
After=rtdd

[Service]
Environment="HOME=/home/root"
#User=root
ExecStart=/usr/sbin/leonard-api
Restart=always
RestartSec=10s
