# leonard/etc/systemd/system/multi-user.target.wants/rtdd.service
[Unit]
Description="Leonard RTD Daemon"
After=zmqd

[Service]
Environment="HOME=/home/root"
#User=root
ExecStart=/usr/sbin/rtdd
Restart=always
RestartSec=10s
